# rootfs/etc/services.d/prometheus-postgres-exporter/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Take down the S6 supervision tree when Prometheus PostgreSQL exporter fails
# ==============================================================================

if [[ "$1" -ne 0 ]] && [[ "$1" -ne 256 ]]; then
  bashio::log.warning "Prometheus PostgreSQL exporter crashed, but continuing add-on"
fi

bashio::log.info "Prometheus PostgreSQL exporter stopped, restarting..." 