# postgresql_postgis/rootfs/etc/services.d/backup-scheduler/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Take down the S6 supervision tree when backup scheduler fails
# ==============================================================================

if [[ "$1" -ne 0 ]] && [[ "$1" -ne 256 ]]; then
  bashio::log.warning "Backup scheduler crashed, but continuing add-on"
fi

bashio::log.info "Backup scheduler stopped, restarting..." 